# claude-managed-service/deploy/my-service.service
[Unit]
Description=My Claude-Managed Task Service
After=network.target

[Service]
Type=simple
# INSTALL_PATH is replaced by the setup skill
WorkingDirectory=INSTALL_PATH
ExecStart=/usr/bin/node INSTALL_PATH/src/index.js
Restart=on-failure
RestartSec=5
Environment=NODE_ENV=production

# Graceful shutdown timeout
TimeoutStopSec=10

[Install]
WantedBy=default.target
